# etc/init/docker-facl.conf
description "run setfacl on /var/run/docker.sock"

start on started docker

task

exec setfacl -m u:huey:rw /var/run/docker.sock

# vim:ft=upstart
