# support/etc/init/huey.conf
description "INLOOP huey workers"

start on starting gunicorn
stop on stopped gunicorn

respawn
respawn limit 3 10

env HOME=/var/lib/inloop

exec envdir /srv/inloop/envdir setuidgid huey django-admin run_huey --workers=2

# vim: ft=upstart
